# rtdb-backup.timer
[Unit]
Description=RethinkDB database off-site backups

[Timer]
OnCalendar=daily
Persistent=true
AccuracySec=30s

[Install]
WantedBy=timers.target
